# pyefis_data/systemd/pyefis-data.service
# systemd USER unit — installs to ~/.config/systemd/user/pyefis-data.service
# Runs one update pass, then refreshes the status JSON the EFIS reads.
#
#   systemctl --user enable --now pyefis-data.timer
#
[Unit]
Description=MakerPlane navigation-data update
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
# A failed/offline update must never be fatal to the aircraft — it just
# means "no new data", which the DATA annunciation surfaces.
# NB: ExecStart path is the pyefis-data entry point. With a venv install,
# symlink it: ln -s ~/pyefis-data-venv/bin/pyefis-data ~/.local/bin/pyefis-data
ExecStart=%h/.local/bin/pyefis-data update
# Refresh the status JSON the EFIS reads (fixed path, not under the data root).
ExecStartPost=-%h/.local/bin/pyefis-data status --json --out %h/.makerplane/pyefis/status.json

[Install]
WantedBy=default.target
